# scripts/minecraft-monitoring.conf
start on runlevel [2345]
stop on runlevel [^2345]

console log

respawn
respawn limit 20 5

exec sudo python /home/ubuntu/MinecraftMonitoring/MineDash/manage.py runserver 0.0.0.0:80
